# deploy/vlog.service
[Unit]
Description=Volo Blog
After=network.target

[Service]
Type=simple
User=vlog
Group=vlog
WorkingDirectory=/opt/vlog
ExecStart=/opt/vlog/server
Restart=on-failure
RestartSec=5
EnvironmentFile=-/etc/vlog/vlog.env
# Default environment if no EnvironmentFile is present
Environment=HOST=127.0.0.1
Environment=PORT=8080
Environment=DATABASE_URL=sqlite:///opt/vlog/storage/vlog.db?mode=rwc
Environment=SESSION_COOKIE_SECURE=1
Environment=RUST_LOG=vlog=info,volo_http=info

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/vlog/storage
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
